# profile/airootfs/etc/systemd/system/ferric-live-user.service
[Unit]
Description=Create FerricOS live user
Before=greetd.service
ConditionPathExists=!/home/liveuser

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ferric-live-user

[Install]
WantedBy=multi-user.target
